# compose-mediaserver.service.template
[Unit]
Requires=docker.service
After=docker.service

[Service]
User=root
Restart=always
ExecStart=/bin/docker-compose -f /usr/local/share/docker/docker-mediaserver/docker-compose.yml up
ExecStop=/bin/docker-compose -f /usr/local/share/docker/docker-mediaserver/docker-compose.yml down

[Install]
WantedBy=default.target
